# config/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
directory=/code
command=blog-service
autostart=true
autorestart=true
